# crontab
SHELL=/bin/bash

# run every day at 03:10
10 3 * * * $HOME/saflii-gazettes/refresh-all.sh >> $HOME/saflii-gazettes/refresh.log 2>&1
